# dotfiles/Documents/PowerShell/profile.d/50-neovim.ps1.tmpl
{{- template "powershell/profile.d/50-neovim.ps1" . -}}
