# .config/bash/profile.d/00_local_opt.bash
LOCAL_OPT="$HOME/.local/opt"

mkdir -p "$LOCAL_OPT/bin"
