# dotfiles/private_dot_config/powershell/profile.d/10-add-local-paths.ps1.tmpl
{{- template "powershell/profile.d/10-add-local-paths.ps1" . -}}
