# dotfiles/private_dot_config/powershell/profile.d/20-environment.ps1.tmpl
chezmoi:template:missing-key=zero
& {{ joinPath .paths.localInstallSbin "set-environment.ps1" }}
