# dotfiles/private_dot_config/powershell/profile.d/31-PSFzf.ps1
Import-Module PSFzf 2>$null
if ($?) {
  # Set-PsFzfOption -PSReadlineChordProvider 'Ctrl+t' -PSReadlineChordReverseHistory 'Ctrl+r'
}
