# dotfiles/private_dot_config/powershell/profile.d/40-prompt.ps1
if (Get-Command oh-my-posh 2>$null) {
    oh-my-posh init pwsh | Invoke-Expression
}
