# dotfiles/private_dot_config/powershell/profile.d/40-prompt.ps1.tmpl
{{- template "powershell/profile.d/40-prompt.ps1" . -}}
