# dotfiles/private_dot_config/powershell/profile.d/50-completions.ps1.tmpl
{{- template "powershell/profile.d/50-completions.ps1" . -}}
