# dotfiles/private_dot_config/powershell/profile.d/50-completions.ps1
foreach ($script in Get-ChildItem $config_dir/powershell/completions.d -Filter '*.ps1') {
    . $script
}
