# dotfiles/private_dot_config/powershell/profile.d/50-neovim.ps1
{{- if eq .chezmoi.os "windows" }}
[System.Environment]::SetEnvironmentVariable('LUNARVIM_CONFIG_DIR', "$HOME/.config/lvim", "Process")
[System.Environment]::SetEnvironmentVariable('LUNARVIM_CONFIG_DIR', "$HOME/.config/lvim", "User")
{{- end }}
