# dotfiles/private_dot_config/powershell/profile.d/zoxide.ps1
if (Get-Command zoxide 2>$null) {
  Invoke-Expression (& { (zoxide init powershell | Out-String) })
}
