# dotfiles/private_dot_config/profile.d/base16.sh.tmpl
{{- if .realConfig.theme | substr 0 7 | eq "base16" -}}
source {{ joinPath .paths.localInstall "base16-shell" "scripts" (printf "%s.sh" .realConfig.theme) }}
{{- end -}}
