# dotfiles/private_dot_config/profile.d/base16.sh.tmpl
{{- if eq "base16" .realConfig.theme -}}
source {{ joinPath .paths.localInstall "base16-shell" "scripts" (printf "base16-%s.sh" .realConfig.base16theme) }}
{{- end -}}
