# dotfiles/private_dot_config/profile.d/base16.sh.tmpl
source {{ joinPath .paths.localInstall "base16-shell" "profile_helper.sh" | quote }}
