# dotfiles/private_dot_config/profile.d/dotnet.sh.tmpl
{{- if eq .chezmoi.os "darwin" -}}
export DOTNET_ROOT=/usr/local/share/dotnet
{{- end -}}
