# .config/powershell/profile.d/50-neovim.ps1##os.Msys
[System.Environment]::SetEnvironmentVariable('LUNARVIM_CONFIG_DIR', "$HOME/.config/lvim", "Process")
[System.Environment]::SetEnvironmentVariable('LUNARVIM_CONFIG_DIR', "$HOME/.config/lvim", "User")
